# examples/com.jaywee.slack-attendance-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jaywee.slack-attendance-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/youruser/Projects/slack-attendance-bot/attendance_bot.py</string>
  </array>

  <!-- Runs daily at 09:05 and 14:05. Add Weekday if you want weekdays only. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>5</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>14</integer>
      <key>Minute</key>
      <integer>5</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/slack-attendance-bot.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/slack-attendance-bot.err</string>
</dict>
</plist>
